# infra/ansible/playbooks/templates/checkpoints.service.j2
[Unit]
Description=Zilliqa Checkpoint Service
After=network.target

[Service]
Type=simple
User=root
Environment=NETWORK_NAME={{ chain_name }}
Environment=ETH_CHAIN_ID={{ eth_chain_id }}
Environment=CHECK_INTERVAL_SECONDS={{ checkpoint_interval_seconds | default('1800') }}
ExecStart=/opt/zilliqa/venv/bin/python3 /usr/local/bin/checkpoints.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target